# evaluation/httpd-2.4.25/nghttp2-1.8.0/contrib/nghttpx.service
[Unit]
Description=HTTP/2 proxy
After=network.target

[Service]
Type=forking
ExecStart=/root/eval/httpd-2.4.25/nghttp2-1.8.0/install/bin/nghttpx --conf=/etc/nghttpx/nghttpx.conf --pid-file=/run/nghttpx.pid --daemon

[Install]
WantedBy=multi-user.target
