# evaluation/httpd-2.4.25/nghttp2-1.8.0/contrib/nghttpx-upstart.conf
# vim: ft=upstart:

description "HTTP/2 reverse proxy"

start on runlevel [2]
stop on runlevel [016]

exec /root/eval/httpd-2.4.25/nghttp2-1.8.0/install/bin/nghttpx
